# scripts/npcd.service.in
[Unit]
Description=Nagios Performance Data C Daemon
After=network.target

[Service]
Type=simple
User=@nagios_user@
Group=@nagios_grp@
ExecStart=@bindir@/@npcd_name@ -f @sysconfdir@/npcd.cfg

[Install]
WantedBy=multi-user.target
